# smf.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!-- 
    Manifest automatically generated by smfgen.
 -->
<service_bundle type="manifest" name="lesbianunix-autoupdate" >
    <service name="lesbianunix/autoupdate" type="service" version="1" >
        <create_default_instance enabled="true" />
        <dependency name="dep0" grouping="require_all" restart_on="error" type="service" >
            <service_fmri value="svc:/milestone/multi-user:default" />
        </dependency>
        <exec_method type="method" name="start" exec="/usr/bin/bash /var/www/lesbianunix.dev/auto &amp;" timeout_seconds="10" >
            <method_context working_directory="/var/www/lesbianunix.dev" >
                <method_credential user="root" group="root" />
                <method_environment >
                    <envvar name="PATH" value="/usr/local/sbin:/usr/local/bin:/opt/local/sbin:/opt/local/bin:/usr/sbin:/usr/bin:/sbin" />
                </method_environment>
            </method_context>
        </exec_method>
        <exec_method type="method" name="stop" exec=":kill" timeout_seconds="30" >
            <method_context working_directory="/var/www/lesbianunix.dev" >
                <method_credential user="root" group="root" />
                <method_environment >
                    <envvar name="PATH" value="/usr/local/sbin:/usr/local/bin:/opt/local/sbin:/opt/local/bin:/usr/sbin:/usr/bin:/sbin" />
                </method_environment>
            </method_context>
        </exec_method>
        <template >
            <common_name >
                <loctext xml:lang="C" >lesbianunix.dev autoupdate script</loctext>
            </common_name>
        </template>
    </service>
</service_bundle>
